# supervisor/supervisor.ini
[program:asset-mgmt]
command=%(ENV_HOME)s/venv/easy-flask/bin/gunicorn --access-logfile - --error-logfile - -b 'unix://%(ENV_HOME)s/socks/easy_flask.sock' src:app -k gevent -w 4 --keep-alive 20
numprocs=1
numprocs_start=0
priority=999
autostart=true
autorestart=unexpected
startsecs=3
startretries=3
exitcodes=0,2
stopwaitsecs=60
directory=%(ENV_HOME)s/project/easy-flask
user=%(ENV_USER)s
stopasgroup=false
killasgroup=false
redirect_stderr=true
stdout_logfile=%(ENV_HOME)s/logs/easy-flask.log
stdout_logfile_maxbytes=250MB
stdout_logfile_backups=10
stderr_logfile=%(ENV_HOME)s/logs/easy-flask.err
stderr_logfile_maxbytes=250MB
stderr_logfile_backups=10
environment=PYTHONPATH='%(ENV_HOME)s/project/easy-flask';LANG=en_US.utf8

